# server.ini
[program:server]
command=python -m SimpleHTTPServer
